# api.joelgrayson.com/crontab
# Every morning
0 0 * * * cd procedures && node save_stats.js

