# api.joelgrayson.com/crontab
# Every morning, save stats
0 0 * * * node procedures/saveStats.js

